# .service
[Unit]
Description=Systemd instance to serve flask App
[Service]
User= "root"
WorkingDirectory= /home/masod_atc
ExecStart= /home/masod_atc/python3 app.py
Restart=always

[Install]
WantedBy=multi-user.target