# .service
[Unit]
Description= Systemd instance to serve flask App

[Service]
User= "flask"
WorkingDirectory= $HOME
ExecStart= $HOME/python3 app.py
Restart=always

[Install]
WantedBy=multi-user.target